# build/apshd.initd
#!/sbin/openrc-run
# Copyright (c) 2022~2024 wangqi <ericwq057@qq.com>
# Distributed under the terms of the MIT License.
#
name="apshd"
description="apshd daemon for the aprilsh service"

command=/usr/bin/apshd
command_args="-verbose 1 2>> /tmp/apshd.log"
command_user=utmp
command_background=yes
directory=/run/utmps
pidfile=/run/utmps/utmpd.pid

depend() {
	need net # declares a hard dependency
}

start_pre() {
	checkpath -d -o utmp:utmp -m 0755 /run/utmps
}
